# overlay.d/30constellation/usr/lib/dracut/modules.d/39constellation-mount/verity-setup.service
[Unit]
Description=Prepare dm-verity mapper device
Before=ostree-prepare-root.service
FailureAction=reboot-immediate

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/sbin/verity-setup
RemainAfterExit=yes
