# overlay.d/30constellation/usr/lib/systemd/system/tpm-pcrs.service
[Unit]
Description=Print PCR state on startup

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/constellation-pcrs

[Install]
WantedBy=multi-user.target

